# ir-bridge.service
[Unit]
Description=IR Remote MQTT Bridge
Wants=network-online.target
After=network-online.target pigpiod.service

[Service]
Type=simple
WorkingDirectory=/var/lib/ir-bridge
EnvironmentFile=/etc/ir-bridge/env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/ir-bridge/venv/bin/python /opt/ir-bridge/mqtt_bridge.py
Restart=always
RestartSec=5
StartLimitIntervalSec=0
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
